# deploy/systemd/moi-souz-api.service
[Unit]
Description=MoiSouz API (NestJS)
After=network.target

[Service]
User=app
WorkingDirectory=/var/www/moi-souz/api
Environment=NODE_ENV=production
Environment=PORT=4000
Environment=HOST=127.0.0.1
ExecStart=/usr/bin/node dist/main.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target